# aso/com.manifestai.growth.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.manifestai.growth</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/yageloshri/ManifestAI-Gratitude-Journal/aso/run_weekly_growth.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>17</integer>
  </dict>
  <key>StandardOutPath</key><string>/Users/yageloshri/ManifestAI-Gratitude-Journal/aso/logs/launchd.out</string>
  <key>StandardErrorPath</key><string>/Users/yageloshri/ManifestAI-Gratitude-Journal/aso/logs/launchd.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
